# deploy/mnemosyne-knowledge-backup.service
[Unit]
Description=Mnemosyne Knowledge Base — daily git snapshot
After=network.target

[Service]
Type=oneshot
# Run as the tree owner (same as the gateway). Once `oste` owns knowledge/.git,
# a root-run `git add/commit` would abort with "detected dubious ownership".
User=oste
Group=oste
WorkingDirectory=/srv/mnemosyne-gb
Environment=KNOWLEDGE_DIR=/srv/mnemosyne-gb/knowledge
ExecStart=/bin/bash /srv/mnemosyne-gb/deploy/knowledge_backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mnemosyne-knowledge-backup
