# deploy/mnemosyne.service
[Unit]
Description=Mnemosyne Knowledge Gateway
Documentation=https://github.com/gborgonovo/mnemosyne-gateway
After=network.target
Wants=network.target

# Stop the restart-loop after 3 failures in 10 min: if the gateway keeps dying
# (e.g. OOM during a mass re-embed), systemd gives up and leaves it in `failed`
# instead of relaunching forever and dragging the whole VPS into swap-thrash.
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple

# Run as the unprivileged owner of the knowledge base, NOT root. Syncthing runs
# as 'oste'; when the gateway ran as root, every atomic rewrite flipped a .md to
# root:root and Syncthing (oste) could no longer read/hash it, so changes stopped
# propagating silently. Running as oste means the gateway and Syncthing share the
# same owner. Requires `oste` to own the whole tree (data/ with the KuzuDB/ChromaDB
# locks, .venv, .env, config/): see deploy/README.md for the one-time chown.
User=oste
Group=oste

WorkingDirectory=/srv/mnemosyne-gb

# Cap memory so a runaway gateway is killed before it starves nginx/SSH. Normal
# steady-state is ~450 MB; a mass re-embed cold boot can spike to ~3 GB and must
# be run manually with the service stopped (see "Re-embed" note below).
MemoryMax=2500M

# Load secrets (OPENAI_API_KEY, MNEMOSYNE_API_KEY, OLLAMA_URL, …)
EnvironmentFile=-/srv/mnemosyne-gb/.env
Environment=PYTHONPATH=/srv/mnemosyne-gb

ExecStart=/srv/mnemosyne-gb/.venv/bin/python3 /srv/mnemosyne-gb/gateway/http_server.py

# Restart automatically on failure. RestartSec keeps systemd from hammering the
# DB if the process exits early; combined with StartLimitBurst above, repeated
# crashes stop the loop instead of looping forever.
Restart=on-failure
RestartSec=60s

# Re-embed note: a deploy that invalidates embeddings (chunking/ID/model change)
# triggers a mass re-embed on next boot that can OOM on a small VPS and loop.
# After such a deploy, run the cold-boot sync manually with the service stopped.
# Run it as 'oste' (NOT root) so it doesn't recreate root-owned files that break
# Syncthing again:
#   sudo systemctl stop mnemosyne
#   cd /srv/mnemosyne-gb && sudo -u oste PYTHONPATH=/srv/mnemosyne-gb \
#     .venv/bin/python3 workers/file_watcher.py --once
#   sudo systemctl start mnemosyne
# (add temporary swap first if RAM is tight).

# Forward stdout/stderr to journald (view with: journalctl -u mnemosyne -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mnemosyne

[Install]
WantedBy=multi-user.target
